# scripts/kiosk-wayland.service
[Unit]
Description=GoldenMunch Kiosk Application (Wayland)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple

# Environment configuration
Environment="NODE_ENV=development"
Environment="ELECTRON_DISABLE_SECURITY_WARNINGS=true"

# CRITICAL: Force X11 and software rendering to prevent DRM/GBM errors
Environment="ELECTRON_OZONE_PLATFORM_HINT=x11"
Environment="ELECTRON_DISABLE_GPU=1"
Environment="LIBGL_ALWAYS_SOFTWARE=1"
Environment="GALLIUM_DRIVER=llvmpipe"
Environment="GBM_BACKEND=dummy"
Environment="MESA_LOADER_DRIVER_OVERRIDE=swrast"
Environment="__GLX_VENDOR_LIBRARY_NAME=mesa"
Environment="DISPLAY=:0"

# Ensure Wayland/X11 display access
PassEnvironment=WAYLAND_DISPLAY XDG_RUNTIME_DIR XDG_SESSION_TYPE DESKTOP_SESSION

# Wait for display to be ready before starting
ExecStartPre=/bin/sleep 5

# Start the kiosk application (uses %h for user home directory)
ExecStart=%h/GoldenMunch_POS-System-With-Custom-Cake-Editor/scripts/start-kiosk-wayland.sh

# Restart policy - retry on failure with backoff
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=120

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=goldenmunch-kiosk

# Resource limits to prevent runaway processes
MemoryLimit=2G
CPUQuota=200%

[Install]
WantedBy=default.target
